# _src/_sbin/partition_db.service
[Unit]
Description=Partition databases

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pitemplog_partition_database.sh